# resources/etc/services.d/dockerd/run
#!/bin/bash
/usr/local/bin/dockerd
echo "dockerd is running"