# resources/etc/services.d/runner/run
#!/usr/bin/with-contenv bash
( \
    ( echo ; echo; echo; echo ) | /runner/config.sh --url "${GIT_URL}" --token "${TOKEN}" \
    && /runner/run.sh \
) \
|| /runner/run.sh